# support/systemd/hics-specim-scanner.service
[Unit]
Description=Hyperspectral Imaging Control System: Specim SWIR scanner control
After=redis-server.service

[Service]
User=fasnachtl
Environment=PYTHONPATH=/home/fasnachtl/git/hics
ExecStart=/usr/bin/python3 -m hics.hardware.specimswir.scanner
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
